# PetaLinux/bsp/uzev/project-spec/meta-avnet/recipes-apps/blinky-systemd/files/blinky.service
[Unit]
Description=Blinky Sample Application
Requires=gpio-utils.service
After=gpio-utils.service

[Service]
Type=simple
ExecStart=/bin/bash /home/root/run-blinky.sh
StandardOutput=tty

[Install]
WantedBy=multi-user.target
