# PetaLinux/bsp/uzev/project-spec/meta-avnet/recipes-utils/gpio-utils-systemd/files/gpio-utils.service
[Unit]
Description=Common GPIO Utilities
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/gpio/gpio_map.sh

[Install]
WantedBy=multi-user.target
